# config/systemd/user/git-sync-obsidian.service
[Unit]
Description=Sync Obsidian Vault with github
Wants=git-sync-obsidian.timer

[Service]
Type=simple
ExecStart=./git/git.sh
