# config/systemd/user/swww.service
[Unit]
Description=Random wallpaper
Wants=swww.timer

[Service]
Type=simple
ExecStart=./swww/swww.sh
