# nix/home-configs/config/systemd/git-sync-obsidian.timer
[Unit]
Description=Run Git Sync for obsidian vault

[Timer]
OnCalendar=*:0/15

[Install]
WantedBy=timers.target
