# nix/home-configs/config/systemd/swww.timer
[Unit]
Description=Run swww wallpaper

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=swww.target
